# Rashma_Docker/php-fpm/crontab
* * * * * root php -q -f /application/artisan schedule:run — no-ansi >> /var/log/cron.log 2>&1
