# crontab
* * * * 15 root php -q -f /application/update
